# roles/graylog_backup/templates/graylog-backup-mongo.timer.j2
[Unit]
Description=Schedule Graylog MongoDB backups

[Timer]
OnCalendar={{ graylog_backup_mongo_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
