# roles/graylog_backup/templates/graylog-datanode-snapshot.timer.j2
[Unit]
Description=Schedule Graylog datanode snapshots

[Timer]
OnCalendar={{ graylog_backup_snapshot_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
