# roles/postgresql_backup/templates/postgresql-backup.timer.j2
[Unit]
Description=Schedule PostgreSQL backups

[Timer]
OnCalendar={{ postgresql_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
