# roles/uptime_kuma_backup/templates/uptime-kuma-backup.timer.j2
[Unit]
Description=Schedule Uptime Kuma backups

[Timer]
OnCalendar={{ uptime_kuma_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
