# roles/plex_backup/templates/plex-backup.timer.j2
[Unit]
Description=Schedule Plex config backups

[Timer]
OnCalendar={{ plex_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
