# roles/invoice_ninja_backup/templates/invoice-ninja-backup.timer.j2
[Unit]
Description=Schedule Invoice Ninja backups

[Timer]
OnCalendar={{ invoice_ninja_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
