# roles/mysql_backup/templates/mysql-backup.timer.j2
[Unit]
Description=Schedule MySQL backups

[Timer]
OnCalendar={{ mysql_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
