# roles/n8n_backup/templates/n8n-backup.timer.j2
[Unit]
Description=Schedule n8n local data backups

[Timer]
OnCalendar={{ n8n_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
